# filesystem/var/lib/dpkg/info/alsa-utils.conffiles
/etc/init/alsa-state.conf
/etc/init/alsa-store.conf
/etc/init/alsa-utils.conf
/etc/init/alsa-restore.conf
/etc/init.d/alsa-utils
