# filesystem/var/lib/dpkg/info/bluez.conffiles
/etc/dbus-1/system.d/bluetooth.conf
/etc/init.d/bluetooth
/etc/bluetooth/rfcomm.conf
/etc/bluetooth/main.conf
/etc/bluetooth/network.conf
/etc/bluetooth/proximity.conf
/etc/bluetooth/input.conf
/etc/bluetooth/audio.conf
/etc/bluetooth/serial.conf
/etc/init/bluetooth.conf
