# filesystem/var/lib/dpkg/info/lightdm.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lightdm.service
/usr
/usr/lib
/usr/lib/lightdm
/usr/lib/lightdm/lightdm-greeter-session
/usr/lib/lightdm/guest-session-auto.sh
/usr/lib/lightdm/lightdm-guest-session
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dm-tool.1.gz
/usr/share/man/man1/lightdm.1.gz
/usr/share/lightdm
/usr/share/lightdm/lightdm.conf.d
/usr/share/lightdm/lightdm.conf.d/50-greeter-wrapper.conf
/usr/share/lightdm/lightdm.conf.d/50-guest-wrapper.conf
/usr/share/lightdm/lightdm.conf.d/50-xserver-command.conf
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_lightdm.py
/usr/share/doc
/usr/share/doc/lightdm
/usr/share/doc/lightdm/copyright
/usr/share/doc/lightdm/changelog.Debian.gz
/usr/share/doc/lightdm/keys.conf
/usr/share/doc/lightdm/lightdm.conf.gz
/usr/sbin
/usr/sbin/lightdm
/usr/sbin/guest-account
/usr/sbin/lightdm-session
/usr/bin
/usr/bin/dm-tool
/etc
/etc/apparmor.d
/etc/apparmor.d/abstractions
/etc/apparmor.d/abstractions/lightdm
/etc/apparmor.d/abstractions/lightdm_chromium-browser
/etc/apparmor.d/lightdm-guest-session
/etc/init.d
/etc/init.d/lightdm
/etc/lightdm
/etc/lightdm/users.conf
/etc/pam.d
/etc/pam.d/lightdm-autologin
/etc/pam.d/lightdm
/etc/pam.d/lightdm-greeter
/etc/X11
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.DisplayManager.conf
/etc/init
/etc/init/lightdm.conf
/var
/var/lib
/var/lib/lightdm-data
/var/cache
/var/cache/lightdm
/var/log
/var/log/lightdm
