# filesystem/var/lib/dpkg/info/ntp.conffiles
/etc/apparmor.d/usr.sbin.ntpd
/etc/apparmor.d/tunables/ntpd
/etc/init.d/ntp
/etc/dhcp/dhclient-exit-hooks.d/ntp
/etc/default/ntp
/etc/cron.daily/ntp
/etc/ntp.conf
