# filesystem/var/lib/dpkg/info/setserial.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/setserial.8.gz
/usr/share/doc
/usr/share/doc/setserial
/usr/share/doc/setserial/README
/usr/share/doc/setserial/pcmcia.repair
/usr/share/doc/setserial/serial.conf
/usr/share/doc/setserial/copyright
/usr/share/doc/setserial/README.Debian.gz
/usr/share/doc/setserial/changelog.Debian.gz
/etc
/etc/init.d
/etc/init.d/etc-setserial
/etc/init.d/setserial
/etc/modutils
/bin
/bin/setserial
/var
/var/lib
/var/lib/setserial
