# filesystem/var/lib/dpkg/info/whoopsie.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/whoopsie.service
/etc
/etc/init
/etc/init/whoopsie.conf
/etc/default
/etc/default/whoopsie
/var
/var/lib
/var/lib/whoopsie
/usr
/usr/bin
/usr/bin/whoopsie
/usr/share
/usr/share/doc
/usr/share/doc/whoopsie
/usr/share/doc/whoopsie/copyright
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/whoopsie.py
/usr/share/doc/whoopsie/changelog.gz
